# club-bot/deploy/club-bot.service
# /etc/systemd/system/club-bot.service
# 鳥人間サークル運営 Discord Bot の常駐サービス定義（さくらのVPS / Ubuntu）
#
# このファイル内のパス・ユーザー名は「作業用ユーザーを ubuntu、
# 配置先を /home/ubuntu/club-bot」とした場合の例です。
# 別のユーザー名・パスにした場合は User / WorkingDirectory /
# ExecStart / EnvironmentFile / Standard* をすべて書き換えてください。
[Unit]
Description=Club Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/club-bot/app
ExecStart=/home/ubuntu/club-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
# .env に DISCORD_TOKEN / ENCRYPTION_KEY / TZ / DB_PATH を設定する。
# 本番（NocoDB 構成）では .env に DATABASE_URL=postgresql://... を追加すると
# PostgreSQL に接続する（未設定なら DB_PATH の SQLite を使用）。
EnvironmentFile=/home/ubuntu/club-bot/.env
StandardOutput=append:/home/ubuntu/club-bot/logs/bot.log
StandardError=append:/home/ubuntu/club-bot/logs/bot.err

[Install]
WantedBy=multi-user.target
